# systemd/atlas-hud.service
# Runs the J.A.R.V.I.S. HUD (Chromium kiosk under cage) on boot.
#
# Install:
#   sudo cp systemd/atlas-hud.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now atlas-hud.service
#
# Before installing, replace YOUR_USERNAME below with your actual Linux
# username (matching atlas-robot.service and atlas-wake.service).
#
# Requires `cage` and `seatd` (see README.md) and the `chromium` package.

[Unit]
Description=A.T.L.A.S. JARVIS HUD Kiosk
After=network-online.target atlas-robot.service
Wants=network-online.target
Requires=atlas-robot.service

[Service]
Type=simple
User=YOUR_USERNAME
Environment=XDG_RUNTIME_DIR=/run/user/1000
# A cold-boot HDMI display can take several seconds to complete its
# handshake after power-on. A blind fixed sleep here isn't long enough on
# some screens — cage fails immediately ("Unable to open Wayland socket")
# if it starts before the DRM connector reports "connected", and it then
# takes several retry cycles (each gated by RestartSec) to eventually
# succeed once the display is actually ready. Poll the connector status
# directly instead of guessing a sleep duration — up to 30s, falling
# through either way so a genuinely stuck/disconnected display doesn't
# block startup forever (Restart=on-failure below is still the safety net).
ExecStartPre=/bin/bash -c 'for i in $(seq 1 30); do grep -q ^connected /sys/class/drm/card*-HDMI-*/status 2>/dev/null && break; sleep 1; done'
ExecStart=/usr/bin/cage -- /usr/bin/chromium --kiosk --no-first-run --disable-restore-session-state --ozone-platform=wayland --start-fullscreen http://127.0.0.1:5051/hud
# There's no physical mouse on this Pi. Chromium can only honor the HUD
# page's `cursor: none` CSS after receiving one real mouse-enter event
# with a Wayland input serial — without a pointer device ever generating
# one, that request never resolves and the default cursor sprite just
# sits on screen forever. This creates a temporary virtual mouse via
# uinput and nudges it once, which is enough to unblock it. Requires the
# udev rule in systemd/99-uinput.rules so the User= account above can
# access /dev/uinput without root.
ExecStartPost=/usr/bin/python3 /home/atlas/atlas-robot/hud_cursor_fix.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
