# systemd/atlas-hud.service.d/10-clean-kiosk.conf
[Service]
# The kiosk runs as a system service, so Chromium does not inherit the
# interactive user's desktop environment. Point it at the real per-user bus
# and Fontconfig file instead of letting it emit misleading startup errors.
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
Environment=FONTCONFIG_FILE=/etc/fonts/fonts.conf

# This display only renders Atlas's local HUD. Disable Chromium background
# services that otherwise try obsolete GCM/component endpoints forever.
ExecStart=
ExecStart=/usr/bin/cage -- /usr/bin/chromium --kiosk --no-first-run --disable-restore-session-state --disable-background-networking --disable-component-update --disable-domain-reliability --disable-sync --disable-features=MediaRouter,OptimizationHints,AutofillServerCommunication --ozone-platform=wayland --start-fullscreen http://127.0.0.1:5051/hud
